# docker/websocket.service
[Unit]
Description=Zoom L20 control WebSocket server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ws
ExecStart=/usr/bin/python3 /opt/ws/websocket_server.py --host 0.0.0.0 --port 8081
Restart=always
RestartSec=3
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=websocket_server
User=root

[Install]
WantedBy=multi-user.target
